# docker/root/etc/services.d/unmanic/run
#!/usr/bin/with-contenv bash

RUN_USER=abc

if [[ -e /app/unmanic/service.py ]]; then
    echo "Update container to running Unmanic from source"
    python_version=$(python3 --version 2>&1 | grep -oP 'Python \K\d+\.\d+')
    if [[ ! -L /usr/local/lib/python${python_version:?}/dist-packages/unmanic ]]; then
        echo "Move container unmanic install"
        mv /usr/local/lib/python${python_version:?}/dist-packages/unmanic /usr/local/lib/python${python_version:?}/dist-packages/unmanic-installed
    fi
    ln -sf /app/unmanic /usr/local/lib/python${python_version:?}/dist-packages/unmanic
fi

unmanic_params=''
if [[ "${DEBUGGING}" == 'true' ]]; then
    unmanic_params="--dev ${unmanic_params}"
fi
if [[ "${USE_TEST_SUPPORT_API}" == 'true' ]]; then
    unmanic_params="--dev-api=https://support-api.test.streamingtech.co.nz ${unmanic_params}"
fi

chmod +x /usr/bin/unmanic
exec s6-setuidgid root /usr/bin/unmanic "${unmanic_params}"
